# wechat-worker-supervisor.conf
[program:wechat-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /app/web/artisan queue:work database --sleep=3 --tries=3 --daemon
autostart=true
autorestart=true
user=www-data
numprocs=8
redirect_stderr=true
stdout_logfile=/app/logs/worker.log